# deploy/systemd/yoro-edge-monitor.service
[Unit]
Description=YORO.gg production edge and recovery monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=yoro
Group=yoro
WorkingDirectory=/srv/yoro/app
EnvironmentFile=/etc/yoro/operations.env
UMask=0077
StateDirectory=yoro-monitor
StateDirectoryMode=0700
ExecStart=/usr/bin/node /srv/yoro/app/scripts/verify-production-edge.mjs --base-url=${PRODUCTION_BASE_URL} --report=/var/lib/yoro-monitor/latest-report.json
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/srv/yoro/backups
