# deploy/systemd/yoro-edge-monitor.timer
[Unit]
Description=Run YORO.gg production edge monitor every five minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
AccuracySec=30s
Persistent=true
Unit=yoro-edge-monitor.service

[Install]
WantedBy=timers.target
